# extractor.service
[Unit]
Description=Extractor Marketing Report Generator
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/extractor
Environment="PATH=/var/www/extractor/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/var/www/extractor"
Environment="FLASK_APP=run.py"
Environment="FLASK_ENV=production"
ExecStart=/var/www/extractor/venv/bin/python /var/www/extractor/venv/bin/gunicorn --config /var/www/extractor/gunicorn_config.py run:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 